# docker/files/etc/supervisor/conf.d/sidekiq.conf
[program:sidekiq]
command=service sidekiq start
