# docker/files/etc/supervisor/conf.d/unicorn.conf
[program:unicorn]
command=service unicorn start
